# scripts/ember-engine-staging.service
[Unit]
After=network.target
Description=Ember crypto copilot service (staging)

[Service]
ExecStartPre=/bin/bash -c 'mv /var/log/ember-engine-stage/output.log /var/log/ember-engine-stage/output-$$(date -Is).log; mv /var/log/ember-engine-stage/error.log /var/log/ember-engine-stage/error-$$(date -Is).log'
ExecStart=/opt/ember-engine/scripts/staging.bash
StandardOutput=file:/var/log/ember-engine-stage/output.log
StandardError=file:/var/log/ember-engine-stage/error.log

[Install]
WantedBy=multi-user.target
